# etc/jpgto.conf
description 'jpgtobot'

start on startup
stop on shutdown

respawn
respawn limit 2 10

script
    chdir /data/jpgtobot
    exec sudo /usr/bin/cabal run
end script
